# gfw-dns/files/gfw-dns.init
#!/bin/sh /etc/rc.common

START=99

. /lib/functions.sh
. /usr/lib/gfw-dns/gfw-dns.sh

start() {
	addrules &
	}

stop() {
	delrules &
	}
